# bakery/repositories/core/recipes/ssh/files/openrc/hydrate-ssh-host-keys
#!/sbin/openrc-run

start() {
    /usr/lib/rugix/scripts/hydrate-ssh-host-keys.sh
}